# swapoff.service
[Unit]
Description=Disable swap

[Service]
Type=simple
User=root
ExecStart=/sbin/swapoff -a
Restart=Always

[Install]
WantedBy=multi-user.target
